# systemd/redshift.service
[Unit]
Description="Start Redshift"

[Service]
ExecStart=/usr/bin/redshift

[Install]
WantedBy=multi-user.target
