# systemd/spotlight.service
[Unit]
Description="Update windows spotlight theme"

[Service]
ExecStart=/etc/theming/scripts/spotlight.sh
User=niek
