# systemd/spotlight.timer
[Unit]
Description=Update Windows spotlight Theme

[Timer]
OnBootSec=5min

[Install]
WantedBy=multi-user.target
